# deploy/shireguard.service
[Unit]
Description=Shireguard WireGuard P2P tunnel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# User= and Environment=HOME= are set by `sudo shireguard install-service`
ExecStart=/usr/local/bin/shireguard up --foreground
Restart=on-failure
RestartSec=5
KillMode=process

# Capabilities needed for TUN device creation (set via setcap during install)
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW

[Install]
WantedBy=multi-user.target
